# ops/dstack/rlab-dstack-image-cleanup.timer
[Unit]
Description=Periodically prune unused rlab dstack training images

[Timer]
OnBootSec=15m
OnUnitActiveSec=30m
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target

